# backend/rc.local
#!/bin/sh -e
. /etc/rc.local.circusd
exit 0